# deploy/vksosa-backend.service
[Unit]
Description=VK SOSA Backend (FastAPI + Uvicorn)
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data

# Путь к проекту — замените на свой (например /var/www/vksosa или /home/user/VK SOSA)
WorkingDirectory=/var/www/vksosa/backend

# Активируйте venv, если используете виртуальное окружение:
# ExecStart=/var/www/vksosa/venv/bin/python -m uvicorn server:app --host 0.0.0.0 --port 8001
ExecStart=/usr/bin/python3 -m uvicorn server:app --host 0.0.0.0 --port 8001

# Переменные окружения (опционально; лучше задать в .env в backend/)
Environment=PYTHONUNBUFFERED=1

Restart=always
RestartSec=5

# Логи
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vksosa-backend

[Install]
WantedBy=multi-user.target
