# deploy/vksosa-frontend.service
[Unit]
Description=VK SOSA Frontend (React, доступ из интернета)
After=network.target vksosa-backend.service
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data

# Путь к проекту — замените на свой
WorkingDirectory=/var/www/vksosa/frontend

# Режим разработки: сайт с WebSocket (HMR), слушает все интерфейсы (0.0.0.0)
# Для продакшена лучше использовать вариант с build + serve (см. DEPLOY_UBUNTU.md)
Environment=HOST=0.0.0.0
Environment=PORT=3000
Environment=REACT_APP_BACKEND_URL=http://127.0.0.1:8001

# yarn или npm — в зависимости от того, что установлено
ExecStart=/usr/bin/yarn start
# Если используете npm: ExecStart=/usr/bin/npm run start

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=vksosa-frontend

[Install]
WantedBy=multi-user.target
